# deploy/systemd/rad.service
[Unit]
Description=rad (rust-analyzer daemon) - user service
After=default.target

[Service]
Type=simple
WorkingDirectory=%h
Environment=RUST_LOG=info
Environment=PATH=%h/.cargo/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/.cargo/bin/rad server
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
