# deploy/systemd/rad.service
[Unit]
Description=rad (rust-analyzer daemon) - user service
After=default.target

[Service]
Type=simple
WorkingDirectory=%h
Environment=RUST_LOG=info
ExecStart=%h/.cargo/bin/rad server
Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=default.target
